# overlay/etc/s6/service/getty/run
#!/usr/bin/env sh

if [ -e /sys/hypervisor/type ] ; then
    console=hvc0
else
    console=ttyS0
fi

exec agetty -l /usr/bin/bash -n -L 38400 $console
#exec agetty -L 38400 $console

