# overlay/etc/s6/service/getty/run
#!/usr/bin/env sh

exec agetty -L 38400 hvc0

