# overlay/etc/s6/service/getty/run
#!/usr/bin/env bash

echo "root:toor" | chpasswd
agetty -L 38400 hvc0

