# overlay/etc/s6/service/getty/run
#!/command/execlineb -P
agetty -L 38400 hvc0

