# overlay/etc/s6/service/getty/run
#!/usr/bin/env sh

if [ "$(cat /sys/hypervisor/type)" == 'xen' ] ; then
    console=hvc0
else
    console=ttyS0
fi

exec agetty -L 38400 $console

