# overlay/etc/s6/service/radvd/run
#!/usr/bin/env sh

exec /usr/bin/radvd --nodaemon --logmethod=stderr -u nobody

