# overlay/etc/s6/service/sshd/run
#!/bin/sh

ssh-keygen -A
exec /usr/bin/sshd -D 2>&1

