# overlay/etc/s6/service/ntpd/run
#!/usr/bin/env sh

exec /usr/bin/ntpd -g -u ntp:ntp -n

