# overlay/etc/s6/service/ntpd/run
#!/usr/bin/env sh

mkdir -p /run/openntpd
exec /usr/bin/ntpd -s -d

